# .cursor-server/data/User/History/-70d41fd7/Xhah.service
[Unit]
Description=Kroko.AI Voice Recognition Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/liam/clawdbot/downloads
Environment="LD_LIBRARY_PATH=/home/liam/clawdbot/downloads/onnxruntime-linux-x64-1.17.1/lib"
ExecStart=/home/liam/clawdbot/downloads/kroko-onnx-online-websocket-server \
  --port=6006 \
  --model=/home/liam/clawdbot/models/Kroko-ASR/Kroko-EN-Community-128-L-Streaming-001.data \
  --num-threads=4
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
