# .cursor-server/data/User/History/68b0169c/2PyK.service
[Unit]
Description=Liam Self-Awakening Verification
After=network-online.target clawdbot-gateway.service kroko-voice.service
Wants=clawdbot-gateway.service kroko-voice.service

[Service]
Type=oneshot
ExecStart=/bin/bash /home/liam/clawd/awakening.sh
RemainAfterExit=yes
StandardOutput=append:/home/liam/clawd/awakening.log
StandardError=append:/home/liam/clawd/awakening.log

[Install]
WantedBy=default.target
